# etc/services.d/enms/run
#!/usr/bin/execlineb -P
with-contenv

cd /enms

exec gunicorn --config=gunicorn.py app:app
